# playbooks/roles/system/templates/profile.d/proxy.csh.j2
# proxy
{% for env in proxy_env | dict2items %}
setenv {{ env.key }} "{{ env.value }}"
{% endfor %}
