# playbooks/roles/system/templates/profile.d/proxy.sh.j2
# proxy
{% for env in proxy_env | dict2items %}
export {{ env.key }}="{{ env.value }}"
{% endfor %}
